# systemd/user/wireplumber.service.d/override.conf
[Service]
ExecStartPost=/home/buggex/.config/wireplumber/bash_scripts/disable_links.sh
